# devops/jraph-app.service
[Unit]
Description=The Jraph web app
[Service]
ExecStart=/bin/bash /root/jraph/devops/run-prod.sh
